# server/quadlet/forge-nginx.container
# Atelier Forge — nginx (HTTP : menus iPXE, WinPE/wimboot, outils de diagnostic).
[Unit]
Description=Atelier Forge nginx (HTTP)
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=forge-nginx
Image=docker.io/library/nginx:stable
Network=host
Volume=@HTTP_ROOT@:/usr/share/nginx/html:ro,z
Volume=@SMB_SHARE@:/srv/deploy:ro,z
Volume=@CONFIG_DIR@/nginx/forge.conf:/etc/nginx/conf.d/default.conf:ro,z

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target default.target
